# systemd/fbcp-ili9341.service
[Unit]
Description=SPI TFT display mirror via fbcp-ili9341
# No After=multi-user.target here: pinboard-kiosk and pinboard-splash are
# WantedBy=multi-user.target *and* After=this unit, so ordering this unit
# after the target closes an ordering cycle and systemd deletes the kiosk
# and splash jobs at boot (issue #67). fbcp only needs /dev/mem and
# restarts on failure, so an early start is safe.

[Service]
Type=simple
# build-safe = rot0 (no DISPLAY_ROTATE_180_DEGREES); build-safe-rot180 is the
# 180°-flipped variant. Same cmake flags otherwise. The display orientation
# and the touch axis inversion in spi_touch_read.c must agree (issue #85).
ExecStart=/home/memomatic/fbcp-ili9341/build-safe/fbcp-ili9341
WorkingDirectory=/home/memomatic/fbcp-ili9341/build-safe
User=root
Restart=always
RestartSec=2
KillSignal=SIGKILL
TimeoutStopSec=2

[Install]
WantedBy=multi-user.target
