# systemd/pinboard-kiosk.service
[Unit]
Description=Memomatic Pinboard kiosk browser
After=pinboard-app.service fbcp-ili9341.service
Wants=pinboard-app.service

[Service]
Type=simple
User=root
# Wait for uvicorn to actually accept connections (Type=simple means 'started'
# as soon as the process spawns, not when it is ready); 2 s is enough on the Pi Zero 2 W.
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/startx /home/memomatic/pinboard/app/kiosk.sh -- :0 -nocursor -s 0 -dpms
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
