# systemd/pinboard-splash.service
[Unit]
Description=Memomatic Pinboard boot splash
# Run after the TFT mirror is up so the image appears on the screen, and
# after getty@tty1 — the kernel console lives on tty1 and agetty clears the
# screen when it starts, which wiped a splash drawn any earlier (issue #70).
After=fbcp-ili9341.service getty@tty1.service
# Finish before the kiosk claims the framebuffer
Before=pinboard-kiosk.service

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/python3 /home/memomatic/pinboard/app/show_splash.py

[Install]
WantedBy=multi-user.target
